# etc/s6-overlay/s6-rc.d/cron/run
#!/command/with-contenv sh
exec cron -f -L 2
